# proxmoxwol-vm.service
[Unit]
Description=Proxmox WOL Listener
After=network.target
Wants=network-online.target
After=network-online.target
After=pve-cluster.service pvedaemon.service
Requires=pve-cluster.service pvedaemon.service

[Service]
Type=simple
Restart=always
RestartSec=5
User=root
ExecStart=/usr/bin/proxmoxwol-listener.py
StandardOutput=journal
StandardError=journal

# Security hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
ReadWritePaths=/var/log/pve

[Install]
WantedBy=multi-user.target
